# layout/DEBIAN/postinst
#!/bin/sh
launchctl load /Library/LaunchDaemons/me.jjolano.shadow.daemon.plist
