# layout/DEBIAN/postinst
#!/bin/sh
launchctl load -w /Library/LaunchDaemons/me.jjolano.shadow.daemon.plist
exit 0
