# layout/DEBIAN/postinst
#!/bin/sh
launchctl load -w /Library/LaunchDaemons/me.jjolano.shadowd.plist
exit 0
